# nexus-oss-carcass/ansible_deploy/roles/nexus_ha_2_nodes/files/nexus-oss-watch.service
[Unit]
Description= Synchronization Script for NEXUS-OSS REPODATA

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/watchers/nexus-oss-watch.sh

[Install]
WantedBy=multi-user.target