# nexus-oss-carcass/ansible_deploy/roles/nexus_ha_2_nodes/files/nexus-oss-watch.path
[Unit]
Description= Tracking NEXUS-OSS REPODATA Files Changes

[Path]
PathModified=/var/lib/docker/volumes/nexus-data
PathModified=/var/lib/docker/volumes/nexus-data/_data

[Install]
WantedBy=multi-user.target